# deploy/systemd/vaulthalla-web.service.in
[Unit]
Description=Vaulthalla Web
After=network.target

[Service]
ExecStart=@BINDIR@/node /usr/share/vaulthalla-web/server.js
Restart=always
User=vaulthalla
Group=vaulthalla
WorkingDirectory=/usr/share/vaulthalla-web
Environment=NODE_ENV=production
Environment=PORT=36968
Environment=HOSTNAME=127.0.0.1
Environment=VAULTHALLA_WEB_INTERNAL_ORIGIN=http://127.0.0.1:36968
Environment=VAULTHALLA_AUTH_ORIGIN=http://127.0.0.1:36970
Environment=VAULTHALLA_PREVIEW_ORIGIN=http://127.0.0.1:36970
Environment=NEXT_PUBLIC_SERVER_ADDR=https://demo.vaulthalla.io

[Install]
WantedBy=multi-user.target
