# deploy/systemd/vaulthalla.service.in
[Unit]
Description=Vaulthalla Service - The Final Cloud
After=network.target
Requires=postgresql.service

[Service]
User=vaulthalla
Group=vaulthalla
SupplementaryGroups=tss
EnvironmentFile=/etc/vaulthalla/vaulthalla.env
ExecStart=@BINDIR@/vaulthalla
Restart=on-failure
RestartSec=3
LimitNOFILE=65535
SyslogIdentifier=vaulthalla
WorkingDirectory=/var/lib/vaulthalla

[Install]
WantedBy=multi-user.target
