# deploy/systemd/vaulthalla-cli.service.in
[Unit]
Description=Vaulthalla control service
After=vaulthalla.service

[Service]
User=vaulthalla
Group=vaulthalla
ExecStart=@BINDIR@/vaulthalla-cli --systemd
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
RestrictAddressFamilies=AF_UNIX
CapabilityBoundingSet=
AmbientCapabilities=
ReadWritePaths=/run/vaulthalla

[Install]
WantedBy=multi-user.target
