# deploy/systemd/vaulthalla-web.service.in
[Unit]
Description=Vaulthalla Web
After=network.target

[Service]
ExecStart=@BINDIR@/node /usr/share/vaulthalla-web/server.js
Restart=always
User=vaulthalla
Group=vaulthalla
WorkingDirectory=/usr/share/vaulthalla-web
Environment=NODE_ENV=production
Environment=PORT=36968
Environment=HOSTNAME=127.0.0.1

[Install]
WantedBy=multi-user.target
