# debian/postinst
#!/bin/sh
set -e

# postinst for vaulthalla
# - Creates system user/group + runtime dirs
# - Performs conservative optional local PostgreSQL bootstrap when available
# - Performs conservative optional nginx integration when available
# - Does NOT prompt routinely during install
# - Does NOT perform package-time super-admin UID binding
#
# Advanced opt-outs:
#   VH_SKIP_DB_BOOTSTRAP=1
#   VH_SKIP_NGINX_CONFIG=1

PKG="vaulthalla"
DB_USER="vaulthalla"
DB_NAME="vaulthalla"
LEGACY_PENDING_SUPERADMIN_UID_FILE="/run/vaulthalla/superadmin_uid"
PENDING_DB_PASS_FILE="/run/vaulthalla/db_password"
CORE_SYSTEMD_UNIT="vaulthalla.service"
CLI_SYSTEMD_UNIT="vaulthalla-cli.service"
CLI_SOCKET_SYSTEMD_UNIT="vaulthalla-cli.socket"
WEB_SYSTEMD_UNIT="vaulthalla-web.service"
PSQL_DEPLOY_DIR="/usr/share/vaulthalla/psql"
NGINX_TEMPLATE="/usr/share/vaulthalla/nginx/vaulthalla"
NGINX_SITE_AVAIL="/etc/nginx/sites-available/vaulthalla"
NGINX_SITE_ENABLED="/etc/nginx/sites-enabled/vaulthalla"
NGINX_MANAGED_MARKER="/var/lib/vaulthalla/nginx_site_managed"

DB_BOOTSTRAP_STATUS="skipped"
NGINX_CONFIG_STATUS="skipped"

create_dir() {
    # create_dir path mode owner group
    install -d -m "$2" -o "$3" -g "$4" "$1"
}

has_command() {
    command -v "$1" >/dev/null 2>&1
}

safe_systemctl() {
    if ! has_command systemctl; then
        return 0
    fi
    systemctl --system "$@" >/dev/null 2>&1 || true
}

run_as_postgres() {
    if has_command runuser; then
        runuser -u postgres -- "$@"
        return $?
    fi
    if has_command sudo; then
        sudo -u postgres "$@"
        return $?
    fi
    return 1
}

env_flag_enabled() {
    case "$1" in
        1|true|TRUE|yes|YES|on|ON)
            return 0
            ;;
        *)
            return 1
            ;;
    esac
}

psql_q() {
    # psql_q "SQL" [database]
    run_as_postgres psql -X -v ON_ERROR_STOP=1 -tAc "$1" "${2:-postgres}"
}

postgresql_server_ready() {
    if has_command pg_isready; then
        run_as_postgres pg_isready -q >/dev/null 2>&1
        return $?
    fi
    run_as_postgres psql -X -v ON_ERROR_STOP=1 -d postgres -tAc 'SELECT 1' >/dev/null 2>&1
}

ensure_tss_group_and_user_added() {
    if getent group tss > /dev/null; then
        echo "[${PKG}] Ensuring '${DB_USER}' is a member of the existing 'tss' group."
    else
        echo "[${PKG}] No 'tss' group found; creating one."
        addgroup --system tss
    fi
    if id -nG "$DB_USER" 2>/dev/null | tr ' ' '\n' | grep -qx tss; then
        return 0
    fi
    adduser "$DB_USER" tss
}

align_config_file_permissions() {
    target="$1"
    if [ ! -f "$target" ]; then
        return 0
    fi
    chown "$DB_USER:$DB_USER" "$target"
    chmod 0640 "$target"
}

require_runtime_assets() {
    if [ ! -d "$PSQL_DEPLOY_DIR" ]; then
        echo "[${PKG}] ERROR: Required SQL deploy directory is missing: ${PSQL_DEPLOY_DIR}" >&2
        exit 1
    fi

    if ! find "$PSQL_DEPLOY_DIR" -maxdepth 1 -type f -name '*.sql' | grep -q .; then
        echo "[${PKG}] ERROR: SQL deploy directory is empty: ${PSQL_DEPLOY_DIR}" >&2
        exit 1
    fi
}

refresh_tpm_udev_rules() {
    if ! has_command udevadm; then
        return 0
    fi
    udevadm control --reload-rules >/dev/null 2>&1 || true
    udevadm trigger --subsystem-match=tpm >/dev/null 2>&1 || true
}

bootstrap_db_if_safe() {
    if env_flag_enabled "${VH_SKIP_DB_BOOTSTRAP:-}"; then
        DB_BOOTSTRAP_STATUS="skipped (VH_SKIP_DB_BOOTSTRAP=1)"
        return 0
    fi

    if ! has_command psql; then
        DB_BOOTSTRAP_STATUS="skipped (psql not installed; install PostgreSQL or configure remote DB)"
        return 0
    fi
    if ! id postgres >/dev/null 2>&1; then
        DB_BOOTSTRAP_STATUS="skipped (local postgres user missing; likely remote/lean install)"
        return 0
    fi
    if ! run_as_postgres true >/dev/null 2>&1; then
        DB_BOOTSTRAP_STATUS="skipped (cannot execute as postgres user)"
        return 0
    fi
    if ! postgresql_server_ready; then
        DB_BOOTSTRAP_STATUS="skipped (local PostgreSQL server not ready)"
        return 0
    fi

    if ! ROLE_EXISTS="$(psql_q "SELECT 1 FROM pg_roles WHERE rolname='${DB_USER}'")"; then
        DB_BOOTSTRAP_STATUS="skipped (failed to query PostgreSQL role state)"
        return 0
    fi
    ROLE_EXISTS="$(printf '%s' "$ROLE_EXISTS" | tr -d '[:space:]')"

    if ! DB_EXISTS="$(psql_q "SELECT 1 FROM pg_database WHERE datname='${DB_NAME}'")"; then
        DB_BOOTSTRAP_STATUS="skipped (failed to query PostgreSQL database state)"
        return 0
    fi
    DB_EXISTS="$(printf '%s' "$DB_EXISTS" | tr -d '[:space:]')"

    ROLE_CREATED="false"
    VAUL_PG_PASS=""

    if [ "$ROLE_EXISTS" != "1" ]; then
        if ! has_command openssl; then
            DB_BOOTSTRAP_STATUS="skipped (openssl not installed for role password generation)"
            return 0
        fi

        VAUL_PG_PASS="$(openssl rand -base64 32 | tr -d '\n')"
        if ! run_as_postgres psql -X -v ON_ERROR_STOP=1 postgres <<EOF
CREATE ROLE ${DB_USER} LOGIN PASSWORD '${VAUL_PG_PASS}';
EOF
        then
            DB_BOOTSTRAP_STATUS="skipped (failed creating PostgreSQL role '${DB_USER}')"
            return 0
        fi
        ROLE_CREATED="true"
        echo "[${PKG}] Created PostgreSQL role '${DB_USER}'."
    fi

    if [ "$DB_EXISTS" != "1" ]; then
        if ! run_as_postgres psql -X -v ON_ERROR_STOP=1 -c "CREATE DATABASE ${DB_NAME} OWNER ${DB_USER};" postgres; then
            DB_BOOTSTRAP_STATUS="skipped (failed creating PostgreSQL database '${DB_NAME}')"
            return 0
        fi
        echo "[${PKG}] Created PostgreSQL database '${DB_NAME}'."
    fi

    if ! run_as_postgres psql -X -v ON_ERROR_STOP=1 -c "GRANT ALL PRIVILEGES ON DATABASE ${DB_NAME} TO ${DB_USER};" postgres; then
        DB_BOOTSTRAP_STATUS="skipped (failed granting DB privileges)"
        return 0
    fi

    if [ "$ROLE_CREATED" = "true" ]; then
        umask 077
        printf '%s\n' "$VAUL_PG_PASS" > "$PENDING_DB_PASS_FILE"
        chown "$DB_USER:$DB_USER" "$PENDING_DB_PASS_FILE"
        DB_BOOTSTRAP_STATUS="done (created local DB role/database and seeded runtime password)"
    else
        DB_BOOTSTRAP_STATUS="done (local DB role/database already present)"
    fi
}

nginx_installed() {
    has_command nginx && [ -d /etc/nginx ]
}

nginx_active() {
    has_command systemctl && systemctl --quiet is-active nginx.service 2>/dev/null
}

has_non_nginx_listener_on_web_ports() {
    if ! has_command ss; then
        return 1
    fi
    WEB_LISTENERS="$(ss -H -ltnp '( sport = :80 or sport = :443 )' 2>/dev/null || true)"
    if [ -z "$WEB_LISTENERS" ]; then
        return 1
    fi
    echo "$WEB_LISTENERS" | grep -v "nginx" | grep -q .
}

has_custom_nginx_sites_enabled() {
    if [ ! -d /etc/nginx/sites-enabled ]; then
        return 1
    fi
    for entry in /etc/nginx/sites-enabled/*; do
        [ -e "$entry" ] || continue
        base="$(basename "$entry")"
        case "$base" in
            default|vaulthalla)
                ;;
            *)
                return 0
                ;;
        esac
    done
    return 1
}

configure_nginx_if_safe() {
    if env_flag_enabled "${VH_SKIP_NGINX_CONFIG:-}"; then
        NGINX_CONFIG_STATUS="skipped (VH_SKIP_NGINX_CONFIG=1)"
        return 0
    fi

    if ! nginx_installed; then
        NGINX_CONFIG_STATUS="skipped (nginx not installed; install nginx for reverse-proxy integration)"
        return 0
    fi
    if ! nginx_active; then
        NGINX_CONFIG_STATUS="skipped (nginx installed but inactive)"
        return 0
    fi
    if has_non_nginx_listener_on_web_ports; then
        NGINX_CONFIG_STATUS="skipped (non-nginx listeners active on 80/443)"
        return 0
    fi
    if has_custom_nginx_sites_enabled; then
        NGINX_CONFIG_STATUS="skipped (custom nginx site layout detected)"
        return 0
    fi
    if [ ! -f "$NGINX_TEMPLATE" ]; then
        NGINX_CONFIG_STATUS="skipped (missing packaged nginx template)"
        return 0
    fi

    install -d -m 0755 /etc/nginx/sites-available /etc/nginx/sites-enabled

    if [ ! -f "$NGINX_SITE_AVAIL" ]; then
        install -m 0644 "$NGINX_TEMPLATE" "$NGINX_SITE_AVAIL"
        printf '%s\n' "managed-by=${PKG}" > "$NGINX_MANAGED_MARKER"
        chmod 0644 "$NGINX_MANAGED_MARKER" >/dev/null 2>&1 || true
    fi

    CREATED_SITE_LINK="false"
    if [ -e "$NGINX_SITE_ENABLED" ] && [ ! -L "$NGINX_SITE_ENABLED" ]; then
        NGINX_CONFIG_STATUS="skipped (${NGINX_SITE_ENABLED} exists and is not a symlink)"
        return 0
    fi

    if [ ! -e "$NGINX_SITE_ENABLED" ]; then
        ln -s "$NGINX_SITE_AVAIL" "$NGINX_SITE_ENABLED"
        CREATED_SITE_LINK="true"
    fi

    if nginx -t >/dev/null 2>&1; then
        safe_systemctl reload nginx.service
        NGINX_CONFIG_STATUS="done (package-managed nginx site enabled)"
        return 0
    fi

    if [ "$CREATED_SITE_LINK" = "true" ]; then
        rm -f "$NGINX_SITE_ENABLED"
    fi
    NGINX_CONFIG_STATUS="skipped (nginx config test failed; reverted site enablement)"
    return 0
}

case "$1" in
  configure)
    if ! id "$DB_USER" >/dev/null 2>&1; then
        adduser --quiet --system --group --home /nonexistent --no-create-home --shell /usr/sbin/nologin "$DB_USER"
    fi

    systemd-tmpfiles --create /usr/lib/tmpfiles.d/vaulthalla.conf || true

    ensure_tss_group_and_user_added
    refresh_tpm_udev_rules

    create_dir /var/lib/vaulthalla 0700 "$DB_USER" "$DB_USER"
    create_dir /var/log/vaulthalla 0750 "$DB_USER" "$DB_USER"
    create_dir /run/vaulthalla     0750 "$DB_USER" "$DB_USER"
    install -d -m 0755 /etc/vaulthalla
    [ -d /mnt/vaulthalla ] || create_dir /mnt/vaulthalla 0755 "$DB_USER" "$DB_USER"
    align_config_file_permissions /etc/vaulthalla/config.yaml
    align_config_file_permissions /etc/vaulthalla/config_template.yaml.in
    require_runtime_assets

    # Remove any stale super-admin seed from older package revisions.
    rm -f "$LEGACY_PENDING_SUPERADMIN_UID_FILE"

    bootstrap_db_if_safe

    safe_systemctl daemon-reexec
    safe_systemctl daemon-reload
    safe_systemctl preset "$CORE_SYSTEMD_UNIT"
    safe_systemctl preset "$CLI_SOCKET_SYSTEMD_UNIT"
    safe_systemctl preset "$CLI_SYSTEMD_UNIT"
    safe_systemctl preset "$WEB_SYSTEMD_UNIT"
    safe_systemctl enable --now "$CORE_SYSTEMD_UNIT"
    safe_systemctl enable --now "$CLI_SOCKET_SYSTEMD_UNIT"
    safe_systemctl enable --now "$CLI_SYSTEMD_UNIT"
    safe_systemctl enable --now "$WEB_SYSTEMD_UNIT"

    configure_nginx_if_safe

    echo "[${PKG}] Install summary:"
    echo "[${PKG}]   DB bootstrap: ${DB_BOOTSTRAP_STATUS}"
    echo "[${PKG}]   nginx config: ${NGINX_CONFIG_STATUS}"
    echo "[${PKG}]   Super-admin ownership: deferred to first CLI use ('vh' or 'vaulthalla')."

    #DEBHELPER#
    ;;

  abort-upgrade|abort-remove|abort-deconfigure)
    #DEBHELPER#
    ;;

  *)
    echo "postinst called with unknown argument \`$1'" >&2
    exit 1
    ;;
esac

exit 0
