# debian/prerm
#!/bin/sh
set -e

PKG="vaulthalla"
WEB_SYSTEMD_UNIT="vaulthalla-web.service"
NGINX_SITE_AVAIL="/etc/nginx/sites-available/vaulthalla"
NGINX_SITE_ENABLED="/etc/nginx/sites-enabled/vaulthalla"
PENDING_SUPERADMIN_UID_FILE="/run/vaulthalla/superadmin_uid"
PENDING_DB_PASS_FILE="/run/vaulthalla/db_password"

has_command() {
    command -v "$1" >/dev/null 2>&1
}

safe_systemctl() {
    if ! has_command systemctl; then
        return 0
    fi
    systemctl --system "$@" >/dev/null 2>&1 || true
}

disable_nginx_site_link_if_package_path() {
    if [ ! -L "$NGINX_SITE_ENABLED" ]; then
        return 0
    fi

    site_target="$(readlink "$NGINX_SITE_ENABLED" 2>/dev/null || true)"
    case "$site_target" in
        "$NGINX_SITE_AVAIL"|../sites-available/vaulthalla)
            rm -f "$NGINX_SITE_ENABLED"
            echo "[${PKG}] Disabled nginx site symlink: ${NGINX_SITE_ENABLED}"
            if has_command nginx && nginx -t >/dev/null 2>&1; then
                safe_systemctl reload nginx.service
            fi
            ;;
        *)
            echo "[${PKG}] Leaving nginx symlink untouched; target is not package site path."
            ;;
    esac
}

case "$1" in
    remove|deconfigure)
        safe_systemctl stop "$WEB_SYSTEMD_UNIT"
        safe_systemctl disable "$WEB_SYSTEMD_UNIT"
        rm -f "$PENDING_SUPERADMIN_UID_FILE" "$PENDING_DB_PASS_FILE"
        disable_nginx_site_link_if_package_path
        ;;
esac

#DEBHELPER#

exit 0
