# debian/prerm
#!/bin/sh
set -e

PKG="vaulthalla"
CORE_SYSTEMD_UNIT="vaulthalla.service"
CLI_SYSTEMD_UNIT="vaulthalla-cli.service"
CLI_SOCKET_SYSTEMD_UNIT="vaulthalla-cli.socket"
WEB_SYSTEMD_UNIT="vaulthalla-web.service"
SWTPM_SYSTEMD_UNIT="vaulthalla-swtpm.service"
NGINX_SITE_AVAIL="/etc/nginx/sites-available/vaulthalla"
NGINX_SITE_ENABLED="/etc/nginx/sites-enabled/vaulthalla"
LEGACY_PENDING_SUPERADMIN_UID_FILE="/run/vaulthalla/superadmin_uid"
PENDING_DB_PASS_FILE="/run/vaulthalla/db_password"
TPM_BACKEND_OVERRIDE_FILE="/etc/systemd/system/vaulthalla.service.d/tpm-backend.conf"

has_command() {
    command -v "$1" >/dev/null 2>&1
}

safe_systemctl() {
    if ! has_command systemctl; then
        return 0
    fi
    systemctl --system "$@" >/dev/null 2>&1 || true
}

disable_nginx_site_link_if_package_path() {
    if [ ! -L "$NGINX_SITE_ENABLED" ]; then
        return 0
    fi

    site_target="$(readlink "$NGINX_SITE_ENABLED" 2>/dev/null || true)"
    case "$site_target" in
        "$NGINX_SITE_AVAIL"|../sites-available/vaulthalla)
            rm -f "$NGINX_SITE_ENABLED"
            echo "[${PKG}] Disabled nginx site symlink: ${NGINX_SITE_ENABLED}"
            if has_command nginx && nginx -t >/dev/null 2>&1; then
                safe_systemctl reload nginx.service
            fi
            ;;
        *)
            echo "[${PKG}] Leaving nginx symlink untouched; target is not package site path."
            ;;
    esac
}

case "$1" in
    remove|deconfigure)
        safe_systemctl stop "$CLI_SOCKET_SYSTEMD_UNIT"
        safe_systemctl stop "$CLI_SYSTEMD_UNIT"
        safe_systemctl stop "$CORE_SYSTEMD_UNIT"
        safe_systemctl stop "$WEB_SYSTEMD_UNIT"
        safe_systemctl stop "$SWTPM_SYSTEMD_UNIT"
        safe_systemctl disable "$CLI_SYSTEMD_UNIT"
        safe_systemctl disable "$CLI_SOCKET_SYSTEMD_UNIT"
        safe_systemctl disable "$CORE_SYSTEMD_UNIT"
        safe_systemctl disable "$WEB_SYSTEMD_UNIT"
        safe_systemctl disable "$SWTPM_SYSTEMD_UNIT"
        rm -f "$LEGACY_PENDING_SUPERADMIN_UID_FILE" "$PENDING_DB_PASS_FILE"
        rm -f "$TPM_BACKEND_OVERRIDE_FILE"
        disable_nginx_site_link_if_package_path
        ;;
esac

#DEBHELPER#

exit 0
